# service_files/active-to-passive.service
[Unit]
Description= Active to Passive

[Service]
ExecStart=/usr/bin/active-to-passive

[Install]
WantedBy=multi-user.target

